# .config/autostart/ranger.desktop
[Desktop Entry]
Name=Ranger
Exec=gnome-terminal --maximize -- ranger
Icon=utilities-terminal
Type=Application